# strato/misc/alpine-cron/rootfs/etc/crontab
# ecr token expires after 12 hours
0 */12 * * * ecr-login.sh | tee -a /var/log/crontab